# nrsysmond.service.d/license.conf
[Service]
Environment="NRSYSMOND_license_key=YOUR_LICENSE_KEY_HERE"
